# config/services/powertop.service
[Unit]
Description=Powertop Tunings

[Service]
Type=oneshot
ExecStart=/usr/bin/powertop --auto-tune

[Install]
WantedBy=multi-user.target
