# scripts/abc.timer
[Unit]
Description=Execute Automatic Brightness Control

[Timer]
OnUnitActiveSec=1m
Unit=abc.service

[Install]
WantedBy=multi-user.target
