# scripts/abc.service
[Unit]
Description=Automatic Brightness Control

[Service]
Type=simple
ExecStart=/usr/bin/abc

[Install]
WantedBy=multi-user.target
